# rootfs/etc/services.d/cron/finish
#!/usr/bin/execlineb -S0

echo "cron exiting"
s6-svscanctl -t /var/run/s6/services
