# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "starting cron"
cron -f -L 15
