# rootfs/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bash

echo "starting influxdb"
influxd run
