# rootfs/etc/services.d/jobber/run
#!/usr/bin/with-contenv bash

echo "starting jobber/cron"
mkdir -p /var/jobber/0/

su -c "/usr/lib/x86_64-linux-gnu/jobberrunner -u /var/jobber/0/cmd.sock /scrutiny/jobber/jobber.yaml" root
