# rootfs/etc/services.d/jobber/run
#!/usr/bin/with-contenv bash

echo "starting jobber/cron"

su -c "/usr/lib/x86_64-linux-gnu/jobberrunner /scrutiny/config/jobber.yaml" root
