# rootfs/etc/services.d/scrutiny/run
#!/usr/bin/with-contenv bash

echo "starting scrutiny"

scrutiny start
